# resticTimer/resticTimer.timer
[Unit]
Description=Custom Backup Restic

[Timer]
; OnBootSec=5min
OnCalendar=*-*-* *:*:00
; OnActiveSec=20s

[Install]
WantedBy=multi-user.target
